# feeders/ping/feeder.ini
[program:soft-watchdog-feeder-ping-__hostname__]
command=__HOME__/.local/bin/soft-watchdog-feeder-ping __hostname__ __interval__
autostart=true
autorestart=true
startretries=99999999
stderr_logfile=/tmp/soft-watchdog-feeder-ping-__hostname__-stderr.log
stdout_logfile=/tmp/soft-watchdog-feeder-ping-__hostname__-stdout.log
user=__whoami__
